# systemd/pacdbserve.socket
[Unit]
Conflicts=pacdbserve.service

[Socket]
ListenStream=7079
Accept=no

[Install]
WantedBy=sockets.target
