# systemd/pacredir.service
[Unit]
Description=Redirect pacman requests via avahi service
Requires=avahi-daemon.service
After=avahi-daemon.service network.target network-online.target
Documentation=https://github.com/eworm-de/pacredir#pacredir

[Service]
Type=notify
ExecStart=/usr/bin/pacredir
DynamicUser=on
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on

[Install]
WantedBy=multi-user.target
Also=pacserve.service
Also=pacdbserve.service
