# systemd/pacserve.service
[Unit]
Description=Darkhttpd to serve pacman cache
After=network.target

[Service]
User=nobody
Group=nobody
ExecStart=/usr/bin/pacserve /var/cache/pacman/pkg/ --ipv6 --port 7078 --no-listing
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on

[Install]
WantedBy=multi-user.target
